# etc/hotplug.d/firewall/97-client_speed_limit
#!/bin/sh

. /lib/functions.sh

[ "$ACTION" = "start" ] && {
    /etc/init.d/client_speed_limit restart
}

