# etc/hotplug.d/firewall/98-guestnetwork_bandwidth_ctrl
#!/bin/sh

. /lib/functions.sh

[ "$ACTION" = "start" ] && {
    /etc/init.d/guestnetwork_bandwidth_ctrl restart
}

