# etc/hotplug.d/firewall/99-qos
#!/bin/sh

if [ "$ACTION" = "add" ]; then
	if [ "$INTERFACE" = "wan" -o "$INTERFACE" = "internet" ]; then
		/etc/init.d/qos restart
	fi
fi
